# packages/nxupdater/finalizer.crontab
*/5 * * * * * * sh /home/git/finalizer.sh