# packages/nxupdater/updater.crontab
*/5 * * * * * * sh /home/node/updater.sh